# deployment/invoice-calculator.service
[Unit]
Description=Invoice Calculator Node.js Application
After=network.target

[Service]
Type=simple
User=www-data
WorkingDirectory=/var/www/invoice-calculator
# Lade Umgebungsvariablen aus .env-Datei (wird automatisch gesetzt)
EnvironmentFile=/var/www/invoice-calculator/.env
Environment="NODE_ENV=production"
Environment="PORT=3000"
ExecStart=/usr/bin/node /var/www/invoice-calculator/dist/index.js
Restart=always
RestartSec=10
StandardOutput=syslog
StandardError=syslog
SyslogIdentifier=invoice-calculator

[Install]
WantedBy=multi-user.target
